# .config/systemd/user/hyprshade_schedule.service
[Unit]
Description=Hyprshade scheduler

[Service]
ExecStart=/home/jamsnxs/.local/bin/hyprshade_schedule
Type=oneshot

